# deploy/systemd/leo-backend.service
[Unit]
Description=Leo ERP Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=leo
Group=leo
WorkingDirectory=/opt/leo/current
Environment=SERVER_PORT=11211
EnvironmentFile=/opt/leo/shared/leo.env
Environment=SPRING_PROFILES_ACTIVE=prod
Environment=LEO_TRUSTED_PROXIES=127.0.0.1
Environment=STEELX_ROOT=/opt/leo
Environment=STEELX_BACKEND_ROOT=/opt/leo
ExecStart=/usr/bin/bash /opt/leo/shared/steelx-process.sh run
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/leo/shared /var/lib/leo /tmp

[Install]
WantedBy=multi-user.target
